# packages/flamenco-manager-go/src/flamenco-manager/flamenco-manager.service
[Unit]
Description=Flamenco Manager
Documentation=https://flamenco.io/
After=mongodb.service

[Service]
Type=simple

ExecStart=/home/flamanager/flamenco-manager
WorkingDirectory=/home/flamanager
User=flamanager
Group=flamanager

Restart=on-failure
RestartSec=1s

EnvironmentFile=-/etc/default/locale

[Install]
WantedBy=multi-user.target
