# packages/flamenco-worker-python/system-integration/sysv/crontab
# Run Flamenco-Worker outside office hours.

SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user  command
0 20 *  * mon-fri root  /usr/sbin/service flamenco-worker start
0  8 *  * mon-fri root  /usr/sbin/service flamenco-worker stop
